# wepo-production-deployment/wepo-node.service.example
[Unit]
Description=WEPO Canonical Full Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wepo
Group=wepo
WorkingDirectory=/opt/wepo
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=WEPO_NETWORK_PROFILE=mainnet
Environment=WEPO_REQUIRE_MAINNET_SEEDS=1
EnvironmentFile=-/etc/wepo/node.env
ExecStart=/opt/wepo/.venv/bin/python /opt/wepo/wepo-blockchain/core/wepo_node.py --network-profile mainnet --data-dir /var/lib/wepo/node --p2p-port 22567 --api-host 127.0.0.1 --api-port 8122 --no-mining
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/wepo /var/log/wepo
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
